# configs/queue-default.container
[Unit]
Description=Laravel queue worker (default)
Requires=laravel-app-network.service laravel-database-network.service storage-volume.service
After=laravel-app-network.service laravel-database-network.service storage-volume.service postgres.service php.service

[Container]
ContainerName=queue-default
Image=localhost/laravel-app:latest
AutoUpdate=local

EnvironmentFile=/etc/laravel/laravel.env
EnvironmentFile=/etc/laravel/php.env
# Workers must not run migrations/autorun (compose: AUTORUN_ENABLED=false).
Environment=AUTORUN_ENABLED=false

Volume=storage.volume:/var/www/html/storage
Network=laravel-app.network
Network=laravel-database.network

Exec=php artisan queue:work --queue=default --tries=3 --backoff=10 --max-time=3600

HealthCmd=healthcheck-queue
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

PodmanArgs=--memory=256m

[Service]
Restart=always
RestartSec=10
StartLimitIntervalSec=0
# compose: stop_grace_period: 60s
TimeoutStopSec=60

[Install]
WantedBy=default.target
